# etc/supervisor.conf

[program:updateSqueue]
command = /usr/local/bin/updateSqueue
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
